# services/forward/dispatcher.service
# Drains the Telegram inbox (captured by telegram-poll), triages each message via the steward and
# routes it to the right agent with a per-agent single-flight queue. Decoupled from the poller so a
# long coach run never blocks pickup. Install to /etc/systemd/system/.
[Unit]
Description=moprox agent dispatcher (Telegram inbox -> steward -> coach/dev/chat)
After=network-online.target telegram-poll.service
Wants=network-online.target
[Service]
Type=simple
User=mikael
Environment=HOME=/home/mikael
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/mikael/projects/moprox-tooling/services/forward/dispatcher.py
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
