# components/dashboard/systemd/dashboard-update.timer
# Run the dashboard updater every 2 minutes (aligned to :00/:02/:04...). Data feeds publish to the
# `data` branch (served via jsDelivr, no Pages build), so a frequent tick keeps hypervisor stats
# fresh without any Pages churn; the updater exits fast when nothing changed. Install to
# /etc/systemd/system/ then `systemctl enable --now dashboard-update.timer` on claude-dev.
[Unit]
Description=Run the moprox dashboard updater every 2 minutes

[Timer]
OnCalendar=*:0/2
Persistent=true

[Install]
WantedBy=timers.target
