# services/mcp-google-standin/mcp-google.service
# STAND-IN: Google Workspace MCP (read-broad, write-narrow) on claude-dev until the mcp-google LXC
# exists. Service account + domain-wide delegation, impersonating one Workspace user; no OAuth2.1
# (incompatible with SA), bound to the trusted subnet.
#
# Permissions use workspace-mcp's granular ladder, reordered by permissions-patch.py so that
# gmail:drafts = readonly + compose (NO modify/trash) and a calendar:events level exists
# (event-write without full calendar ACL). Re-run permissions-patch.py after any pip upgrade.
#
# DWD scope authorization (admin.google.com) gates what the SA may request — the union of all
# enabled scopes is requested in ONE token, so an unauthorized scope breaks the WHOLE server.
# Authorized so far: gmail.readonly, gmail.compose, calendar.readonly, drive.readonly,
# documents.readonly. calendar.events is the one outstanding admin step; once added, flip
# calendar:readonly -> calendar:events below (target end-state) and restart.
[Unit]
Description=Google Workspace MCP (workspace-mcp) — STAND-IN on claude-dev (read + gmail drafts)
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
User=mikael
Environment=HOME=/home/mikael
Environment=GOOGLE_SERVICE_ACCOUNT_KEY_FILE=/home/mikael/.config/claude-dev/google-sa.json
Environment=USER_GOOGLE_EMAIL=mikael@odinlake.net
Environment="WORKSPACE_MCP_PERMISSIONS=gmail:drafts calendar:readonly drive:readonly docs:readonly"
Environment=WORKSPACE_MCP_HOST=0.0.0.0
Environment=WORKSPACE_MCP_PORT=8000
WorkingDirectory=/opt/mcp-google
ExecStart=/opt/mcp-google/venv/bin/workspace-mcp --transport streamable-http
Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target
